# transcoder-pipeline/systemd/trigger-organize.path
[Unit]
Description=Watch elite trigger for manual organize run

[Path]
PathExists=/var/lib/elite-triggers/organize

[Install]
WantedBy=multi-user.target
