# transcoder-pipeline/systemd/trigger-transcode.path
[Unit]
Description=Watch elite trigger for manual transcode run

[Path]
PathExists=/var/lib/elite-triggers/transcode

[Install]
WantedBy=multi-user.target
